# targets/debian/postinst
#!/bin/sh

NAME=hakit

update-rc.d $NAME defaults
update-rc.d $NAME enable
service $NAME stop
service $NAME start
