# targets/debian/prerm
#!/bin/sh

NAME=hakit

service $NAME stop
update-rc.d -f $NAME remove
